# build/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:sshd]
command=/usr/sbin/sshd -D

[program:rsyncd]
command=/usr/bin/rsync --no-detach --daemon --log-file=/dev/stdout